# design/systemd/bitmover/bitmover.service
# /usr/lib/systemd/system/bitmover.service

[Unit]
Description=Data Mover Application
# Documentation path should point to where your RPM installs it (DM-NFR-004)
Documentation=file:///usr/share/doc/datamover-app/README.md
# Ensures Data Mover starts after network and the exportcliv8 services group (DM-FR-START-003)
After=network.target exportcliv8.target
# If bitmover.service is started, systemd will also try to start exportcliv8.target
Wants=exportcliv8.target

# Crash loop protection: stop trying if it fails too often (DM-NFR-001)
# Example: no more than 5 restarts in 5 minutes (300 seconds)
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
# 'simple' is suitable as your Python app's main_entrypoint runs in the foreground
# and handles signals for graceful shutdown.
Type=simple

User=datamover_user
Group=datapipeline_group
UMask=0027

# ExecStart: Command to start the Data Mover application.
# The path to 'datamover_executable' will depend on how your Python wheel's
# entry_points (e.g., console_scripts) are defined and where the RPM installs it.
# It might be in a virtual environment's bin or a system bin.
# The '--config' argument path should also match your deployment. (DM-GEN-003, DM-FR-START-002)
ExecStart=/usr/bin/datamover_executable --config /etc/datamover_app/app.ini
#
# Example if installed in a venv at /opt/datamover_app/venv:
# ExecStart=/opt/datamover_app/venv/bin/datamover_executable --config /etc/datamover_app/app.ini
#
# Or if directly calling the python script (less common for wheel-installed apps with entry points):
# ExecStart=/usr/bin/python3 /opt/datamover_app/main_script.py --config /etc/datamover_app/app.ini

# Set PYTHONUNBUFFERED=1 to ensure Python's output (especially logging)
# goes directly to journald without extra buffering, which is good for services.
Environment="PYTHONUNBUFFERED=1"

# Optional: Set a working directory if your application expects to run from a specific path.
# For example, if it uses relative paths for resources not covered by the config.
# WorkingDirectory=/opt/datamover_app/
# Or, if it primarily operates on /opt/bit_mover:
# WorkingDirectory=/opt/bit_mover/

# Restart behavior (DM-FR-RESTART-001)
Restart=on-failure
# Time to wait before attempting a restart
RestartSec=10s

# Prevent restarts on specific fatal error codes from sysexits.h (DM-FR-EXITCODES-001)
# Your application code confirmed it uses:
# EX_USAGE (64) for command line usage error
# EX_CONFIG (78) for configuration error
# AppSetupError within run() also exits with EX_CONFIG.
RestartPreventExitStatus=64 78

# Logging configuration: send stdout/stderr to journald (DM-NFR-002)
StandardOutput=journal
StandardError=journal
# Tag log entries in the journal with 'datamover' (DM-NFR-002)
SyslogIdentifier=datamover

[Install]
# Ensures the service is enabled to start automatically on system boot (DM-FR-START-001)
WantedBy=multi-user.target